# meta-eris-linux/recipes-eris/eris-rest-api/eris-rest-api/eris-rest-api.service
[Unit]
Description=Eris Linux REST API

[Service]
ExecStart=/usr/sbin/eris-rest-api

[Install]
WantedBy=multi-user.target
