# meta-eris-linux/recipes-eris/eris-ntp/eris-ntp/eris-ntp.service
[Unit]
Description=Eris Linux NTP service

[Service]
ExecStart=/usr/sbin/eris-ntp

[Install]
WantedBy=multi-user.target
