# meta-eris-linux/recipes-eris/eris-api-server/eris-api-server/eris-api-server.service
[Unit]
Description=Eris Linux API server

[Service]
ExecStart=/usr/sbin/eris-api-server

[Install]
WantedBy=multi-user.target
